# get-dns-blocklists.timer
[Unit]
Description=Update DNS blacklists daily

[Timer]
OnCalendar=daily
Persistent=false

[Install]
WantedBy=timers.target
